# 7.1/rootfs/etc/supervisor/conf.d/app-log.conf
[program:app-log]
command=/bin/bash -c "tail -q -n 0 -F app/logs/dev.log app/logs/prod.log var/logs/dev.log var/logs/prod.log"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
