# 7.2/rootfs/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm7.2 --nodaemonize
