# config/supervisord/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php5-fpm --nodaemonize
